# infra/backup/com.kbo.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- task-001 — Supabase pg_dump 일일 백업 LaunchAgent (cron 대체: macOS FDA 이슈 회피)
     설치: cp 이 파일 ~/Library/LaunchAgents/ && launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.kbo.backup.plist
     매일 04:00 실행. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kbo.backup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/arden/code/kbo-away-fans/scripts/backup/pg_dump.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/arden/Library/Logs/kbo-backup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/arden/Library/Logs/kbo-backup.log</string>
</dict>
</plist>
